# change-wallper.service
[Service]
Type=oneshot
ExecStart=/bin/bash -c '~/.local/bin/change-wallper.sh'
